# docker/php-fpm/crontab
10 22 * * * php /app/bin/console app:AddWaitingImplementationToTesting >> /var/log/schedule.log 2>&1